# initramfs/etc/inittab
::sysinit:/bin/mount -a
::sysinit:/bin/hostname -F /etc/hostname
::respawn:/bin/cttyhack /bin/ash
